# tests/integration/targets/setup_letsencrypt/templates/pebble-challtestsrv.service.j2
[Service]
ExecStart=/usr/local/bin/pebble-challtestsrv -defaultIPv6 ""
Restart=always
StandardOutput=journal
User=root
Group=root

[Install]
WantedBy=multi-user.target
