# tests/integration/targets/setup_letsencrypt/templates/pebble-challtestsrv.service.j2
[Service]
ExecStart=/usr/local/bin/pebble-challtestsrv
Restart=always
StandardOutput=syslog
User=root
Group=root

[Install]
WantedBy=multi-user.target
