# tests/integration/targets/setup_letsencrypt/templates/pebble.service.j2
[Service]
ExecStart=/usr/local/bin/pebble -dnsserver localhost:8053
Restart=always
StandardOutput=syslog
User=root
Group=root

[Install]
WantedBy=multi-user.target
